# systemd/deep-research-nightly-verification.timer
[Unit]
Description=Run Deep Research knowledge-base verification nightly

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=deep-research-nightly-verification.service

[Install]
WantedBy=timers.target
